# contrib/ailm.service
[Unit]
Description=ailm — AI Linux Monitor
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/ailm --no-ui
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
